# scripts/launchctl/com.pinkybot.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pinkybot.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/oleg/PinkyBot &amp;&amp; set -a &amp;&amp; source .env &amp;&amp; set +a &amp;&amp; exec .venv/bin/python -m pinky_daemon --mode api --port 8888 --host 0.0.0.0 --working-dir .</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/oleg/PinkyBot</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/oleg/PinkyBot/logs/api.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/oleg/PinkyBot/logs/api.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
